# deploy/guru-web.service
[Unit]
Description=Guru web app (Next.js)
After=network.target postgresql.service ollama.service
Requires=postgresql.service ollama.service

[Service]
Type=simple
User=guru
Group=guru
WorkingDirectory=/srv/guru-web/current
EnvironmentFile=/etc/guru-web.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

# Hardening
# NOTE: MemoryDenyWriteExecute=true is intentionally OMITTED. V8's JIT compiler
# requires writable+executable memory pages (mprotect with PROT_WRITE|PROT_EXEC),
# which that directive blocks. With it set, Node panics at startup with
# "Check failed: 12 == errno" inside V8's MemoryAllocator. The other hardening
# directives below are JIT-safe.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/guru-web/current
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
